# defaultConfig/etc/skel/.config/autostart/xcalib.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=xcalib
Exec=~/.local/scripts/xcalibProfile.sh
StartupNotify=false
Terminal=false
