# defaultConfig/etc/profile.d/kwin.sh
#!/bin/sh
# if there is screen tearing uncomment
#export KWIN_TRIPLE_BUFFER=1 
